# docker/conf/supervisord.conf
[supervisord]
nodaemon=true
identifier=supervisor
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
umask=002
user=www-data

[unix_http_server]
file=/tmp/supervisor.sock
chown=www-data:www-data
chmod=0770

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:php-fpm]
command=php-fpm -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0